# vagrant/st2/rc.local
#!/bin/sh -e

ssh-agent -s
ssh-add /tmp/vagrant_key
